# deploy/systemd/okx-dashboard-tunnel.service
[Unit]
Description=OKX Quant Dashboard Cloudflare Tunnel
After=network-online.target okx-dashboard-proxy.service
Wants=network-online.target
Requires=okx-dashboard-proxy.service

[Service]
Type=simple
User=okxbot
Group=okxbot
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --protocol http2 --url http://127.0.0.1:8766
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
